# launchd/com.skufs.tutor-backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skufs.tutor-backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/andriy/tutor_skufs/backend/.venv/bin/python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>4800</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/andriy/tutor_skufs/backend</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/andriy/tutor_skufs/backend/logs/out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andriy/tutor_skufs/backend/logs/err.log</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
